# cloud-function/scripts/supervisor_template.conf
# 📦 Supervisor Configuration Template for Trading Bot
# Copy this to ~/supervisor/trading_bot.conf and adjust paths

[supervisord]
logfile=%(ENV_HOME)s/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=%(ENV_HOME)s/supervisor/supervisord.pid
nodaemon=false

[supervisorctl]
serverurl=unix:///%(ENV_HOME)s/supervisor/supervisor.sock

[unix_http_server]
file=%(ENV_HOME)s/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:trading_bot]
# Command to run (ADJUST PYTHON PATH if needed)
command=/usr/local/bin/python3 %(ENV_HOME)s/code/stockScreener/cloud-function/scripts/trading_bot.py

# Working directory
directory=%(ENV_HOME)s/code/stockScreener/cloud-function

# Auto-start on supervisord startup
autostart=true

# Auto-restart if crashes
autorestart=true
autorestart_retries=10
exitcodes=0,2

# Wait 30 seconds before sending SIGKILL after SIGTERM
stopwaitsecs=30

# Redirect logs
stderr_logfile=%(ENV_HOME)s/supervisor/trading_bot.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=%(ENV_HOME)s/supervisor/trading_bot.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# User to run as
user=%(ENV_USER)s

# Environment variables
environment=TRADING_ENVIRONMENT="demo"

# Priority (lower starts first)
priority=999

# Don't restart too quickly
startsecs=10
startretries=3

# Capture log rotation signals
redirect_stderr=true
